# ylx_files/x2/files/etc/hotplug.d/iface/00-netstate
[ ifup = "$ACTION" ] && {
	uci_toggle_state network "$INTERFACE" up 1
	[ -n "$DEVICE" ] && {
		uci_toggle_state network "$INTERFACE" ifname "$DEVICE"
	}
    [ "uci get shadowsocksr.@global[0].global_server" != "nil"  ] && {
    /etc/init.d/shadowsocksr restart &
    }
}
